# boot_as_service/nostalgia.service
[Unit]
Description=nostalgia service

[Service]
Type=simple
ExecStart=/usr/bin/nostalgia serve

[Install]
WantedBy=multi-user.target
